# shim/systemd/chariot.service
[Service]
Environment="RUST_LOG=debug"
# EnvironmentFile=
ExecStart=/usr/bin/chariot-shim --config /etc/chariot/shim.toml